# config/dyn-pxe-watcher.path
[Path]
Unit=dyn-pxe-watcher.service
PathChanged=/etc/dyn-pxe

[Install]
WantedBy=multi-user.target
